# deploy/linux/uncrumpled.socket
# Uncrumpled Context Switcher - systemd Socket Activation
#
# This file enables socket activation for the Uncrumpled daemon.
# The daemon will be started on-demand when a client connects.
#
# Installation:
#   mkdir -p ~/.config/systemd/user
#   cp uncrumpled.socket ~/.config/systemd/user/
#   cp uncrumpled@.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable uncrumpled.socket
#   systemctl --user start uncrumpled.socket
#
# Test:
#   nc -U $XDG_RUNTIME_DIR/uncrumpled/uncrumpled.sock

[Unit]
Description=Uncrumpled Context Switcher Socket
Documentation=https://github.com/uncrumpled/context-switcher

[Socket]
ListenStream=%t/uncrumpled/uncrumpled.sock
SocketMode=0600

# Remove stale socket file on start
RemoveOnStop=yes

# Accept connections in the daemon
Accept=no

[Install]
WantedBy=sockets.target
