# deploy/linux/uncrumpled-context-switcher-socket.service
# Uncrumpled Context Switcher - Socket-Activated Service
#
# This service variant is started via socket activation.
# The systemd socket unit (uncrumpled-context-switcher.socket) passes the socket FD.
#
# See uncrumpled-context-switcher.socket for installation instructions.

[Unit]
Description=Uncrumpled Context Switcher Daemon (Socket Activated)
Documentation=https://github.com/uncrumpled/uncrumpled-context-switcher
Requires=uncrumpled-context-switcher.socket
After=uncrumpled-context-switcher.socket

[Service]
Type=simple
ExecStart=%h/.local/bin/uncrumpled-context-switcher-daemon --foreground --systemd-socket
ExecReload=/bin/kill -HUP $MAINPID

# Accept socket from systemd
StandardInput=socket

# Use runtime directory from socket unit
RuntimeDirectory=uncrumpled-context-switcher
RuntimeDirectoryMode=0700

# Resource limits
LimitNOFILE=1024
MemoryMax=128M

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadWritePaths=%h/.config/uncrumpled-context-switcher

[Install]
WantedBy=default.target
