# deploy/linux/uncrumpled.service
# Uncrumpled Context Switcher - systemd User Service
#
# This is a systemd user service file for the Uncrumpled daemon.
# It should be installed to ~/.config/systemd/user/uncrumpled.service
#
# Installation:
#   mkdir -p ~/.config/systemd/user
#   cp uncrumpled.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable uncrumpled
#   systemctl --user start uncrumpled
#
# Useful commands:
#   systemctl --user status uncrumpled
#   systemctl --user restart uncrumpled
#   journalctl --user -u uncrumpled -f   # View logs

[Unit]
Description=Uncrumpled Context Switcher Daemon
Documentation=https://github.com/uncrumpled/context-switcher
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/uncrumpled-daemon --foreground
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

# Socket and state file paths
Environment=XDG_RUNTIME_DIR=%t
RuntimeDirectory=uncrumpled
RuntimeDirectoryMode=0700

# Resource limits
LimitNOFILE=1024
MemoryMax=128M

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Allow access to user config directory
ReadWritePaths=%h/.config/uncrumpled

[Install]
WantedBy=default.target
